# jobhuntwow-app/agent/supervisord.conf
# Reused from cybergodai/cybergodai. Runs the virtual-display + VNC + Chrome stack.
# The agent itself is NOT in supervisord — it's triggered on demand via `python jhw.py scrape/apply`
# (docker exec), so it connects to the already-running Chrome and the candidate watches via noVNC.
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:xvfb]
command=/usr/bin/Xvfb :0 -screen 0 1280x1024x24 -ac -nocursor
user=chrome-user
autorestart=true
priority=1
startsecs=2
startretries=20

[program:x11vnc]
command=/usr/bin/x11vnc -display :0 -forever -shared -nopw -listen 0.0.0.0 -xkb
user=chrome-user
autorestart=true
priority=2
startsecs=2
startretries=20

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc 9090 0.0.0.0:5900
user=chrome-user
autorestart=true
priority=3
startsecs=2
startretries=20

[program:chrome]
command=/usr/local/bin/run-chrome.sh
environment=DISPLAY=":0",HOME="/home/chrome-user",XDG_RUNTIME_DIR="/tmp/runtime-chrome-user"
user=chrome-user
autorestart=true
priority=10
startsecs=10
startretries=10
